# root/etc/services/nginx/run
#!/bin/sh -e

generate-self-signed-ssl
exec nginx
