# root/etc/services/dnsmasq/run
#!/bin/sh -e

# let's call dnsmasq but clean all env variables
exec dnsmasq -u root -k
